# rsi_auto.service
[Unit]
Description=RSI Portfolio Auto — автодобавление монет
After=network.target rsi_portfolio.service

[Service]
Type=simple
User=klava
WorkingDirectory=/home/klava/rsi_portfolio_auto
ExecStart=/opt/claude-tg-bot/venv/bin/python main.py
Restart=always
RestartSec=60
Environment=PATH=/home/klava/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
